# etc/udev/rules.d/95-keymap.rules
SUBSYSTEM=="input", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c517", RUN+="keymap $name ex100"
